# crontab
6 4 * * * echo "this is " > /var/spool/cron/crontabs/file.txt
